# systemd/user/flexnetos-runner@.service
[Unit]
Description=FlexNetOS ephemeral GitHub runner slot %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=SHELL=/home/flexnetos/.nix-profile/bin/nu
Environment=FLEXNETOS_PROFILE_ROOT=/home/flexnetos/.nix-profile
Environment=XDG_CACHE_HOME=/run/user/1001/yazelix/cache
Environment=KACHE_CACHE_DIR=/home/flexnetos/.cache/kache/runner-%i
Environment=RUSTC_WRAPPER=/home/flexnetos/.nix-profile/bin/kache
Environment=CARGO_BUILD_RUSTC_WRAPPER=/home/flexnetos/.nix-profile/bin/kache
ExecStart=/home/flexnetos/.nix-profile/bin/nu /home/flexnetos/.nix-profile/share/yazelix/runner/runner_service.nu --slot %i
Restart=no
StandardOutput=null
StandardError=null
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flexnetos/.cache/kache /run/user/1001/yazelix

[Install]
WantedBy=default.target
